# etc/services.d/healthcheck/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

sleep 900
exec \
    /scripts/health_from_adsbexchange.py \
    | awk -W interactive '{print "[healthcheck] " $0}'
