# etc/services.d/adsbexchange-stats/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    /usr/local/bin/json-status 2>&1 | awk -W interactive '{print "[adsbexchange-stats] " $0}'
